# systemd/osg-cleanup.timer
[Unit]
Description=Run OSG Cleanup hourly, 15 minutes after the hour

[Timer]
OnCalendar=*-*-* *:15:00

