# cron.d/osg-cleanup
# Cron job running by default every 1 hour.
# The lock file can be enabled or disabled via a
# service osg-cleanup start
# chkconfig osg-cleanup on

# Note the lock file not existing is success hence the the slightly odd logic
# below.

15 1 * * *    root    [ ! -f /var/lock/subsys/osg-cleanup-cron ] || /usr/sbin/osg-cleanup
